# systemd/crm3-drive-backup.service
[Unit]
Description=CRM 3.0 Google Drive backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/crm3
EnvironmentFile=/root/crm3/.env
Environment=GOOGLE_DRIVE_AUTH_MODE=oauth
Environment=GOOGLE_DRIVE_OAUTH_REFRESH_TOKEN_FILE=/root/crm3/runtime-secrets/google_drive_oauth_refresh_token
Environment=BACKUP_DB_FALLBACK_HOST=localhost
Environment=GOOGLE_DRIVE_BACKUP_FOLDER_ID=1UEpbPVnt9wZ-y56GCq_GBOurBJHCNoVf
Environment=GOOGLE_DRIVE_ROOT_FOLDER_ID=1t9-ZijytjS5BE8Kb_HDTpwJrEsP4Pc3Z
ExecStart=/root/crm3/backend/.venv/bin/python scripts/backup_project_to_drive.py
Nice=10
StandardOutput=append:/root/crm3/cron-backup.log
StandardError=inherit
