# systemd/crm3-drive-backup.service
[Unit]
Description=CRM 3.0 Google Drive backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/crm3
ExecStart=/bin/bash -lc 'BACKUP_DB_FALLBACK_HOST=localhost GOOGLE_DRIVE_SERVICE_ACCOUNT_FILE=/root/crm3/credentials.json GOOGLE_DRIVE_BACKUP_FOLDER_ID=1UEpbPVnt9wZ-y56GCq_GBOurBJHCNoVf GOOGLE_DRIVE_ROOT_FOLDER_ID=1t9-ZijytjS5BE8Kb_HDTpwJrEsP4Pc3Z /root/crm3/backend/.venv/bin/python scripts/backup_project_to_drive.py >> /root/crm3/cron-backup.log 2>&1'
Nice=10
